# share/apache/fcgiodata.socket
[Unit]
Description=fast cgi socket for odata

[Socket]
ListenStream=127.0.0.1:4200
Accept=false

[Install]
WantedBy=multi-user.target network.target
