# share/apache/fcgi-threaded-hello-world.socket
[Unit]
Description=fast cgi socket for odata

[Socket]
ListenStream=127.0.0.1:5200
Accept=false

[Install]
WantedBy=sockets.target
